# ansible/playbooks/roles/restic/templates/etc/systemd/system/restic.timer.j2
# Managed by Ansible

[Unit]
Description=Run Restic on a schedule

[Timer]
OnCalendar={{ restic_schedule }}

[Install]
WantedBy=timers.target
